# scripts/com.olsentribune.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">

<!--
  com.olsentribune.daily.plist
  ─────────────────────────────────────────────────────────────────────────────
  This is a macOS launchd "agent" — a background job that macOS runs on a
  schedule automatically, even if you're not logged in at that moment.

  HOW TO INSTALL:
    1. Copy this file to ~/Library/LaunchAgents/
         cp com.olsentribune.daily.plist ~/Library/LaunchAgents/

    2. Tell macOS to load it (registers the schedule):
         launchctl load ~/Library/LaunchAgents/com.olsentribune.daily.plist

  HOW TO TEST IT RIGHT NOW (without waiting for 8 AM):
         launchctl start com.olsentribune.daily

  HOW TO UNINSTALL:
         launchctl unload ~/Library/LaunchAgents/com.olsentribune.daily.plist
         rm ~/Library/LaunchAgents/com.olsentribune.daily.plist

  LOGS:
    - /tmp/olsen-tribune-daily.log  — script output (articles generated, errors)
    - /tmp/olsen-tribune-daily.err  — any unexpected OS-level launch errors
  ─────────────────────────────────────────────────────────────────────────────
-->

<plist version="1.0">
<dict>

  <!-- Unique identifier for this job — must match the filename (without .plist) -->
  <key>Label</key>
  <string>com.olsentribune.daily</string>

  <!-- The shell script to run — uses /bin/bash so that the PATH extension
       inside run-daily.sh works correctly -->
  <key>ProgramArguments</key>
  <array>
    <!-- Call node directly — avoids macOS security restrictions on shell scripts -->
    <string>/usr/local/bin/node</string>
    <string>/Users/chrisolsen/Desktop/the-olsen-tribune/scripts/daily-update.js</string>
  </array>

  <!-- Set the working directory so relative paths inside the script resolve correctly -->
  <key>WorkingDirectory</key>
  <string>/Users/chrisolsen/Desktop/the-olsen-tribune/scripts</string>

  <!-- Schedule: run every day at 08:00 local time.
       launchd uses your Mac's system clock timezone.
       If your Mac is set to Eastern Time, this fires at 8:00 AM ET.
       If you're in a different timezone, adjust Hour accordingly:
         Eastern  (ET/EDT) →  8
         Central  (CT/CDT) →  8  (set Mac clock to CT)
         Pacific  (PT/PDT) →  8  (set Mac clock to PT) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- If the Mac was asleep at 8 AM, run the job as soon as it wakes up -->
  <key>RunAtLoad</key>
  <false/>

  <!-- Where to write normal log output (script progress, articles generated) -->
  <key>StandardOutPath</key>
  <string>/tmp/olsen-tribune-daily.log</string>

  <!-- Where to write error output (crashes, missing files, etc.) -->
  <key>StandardErrorPath</key>
  <string>/tmp/olsen-tribune-daily.err</string>

</dict>
</plist>
